# mailcatcher/mailcatcher.ini
[program:mailcatcher]
command=mailcatcher --foreground --ip=0.0.0.0 --smtp-port=25 --http-port=88
autostart=true
autorestart=true
priority=5
startretries=3
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
